# foreman/hiringapp-worker.conf
start on starting hiringapp
stop on stopping hiringapp
